# agent-workspace/deploy/esnafdijital-admin.service
[Unit]
Description=EsnafDigital admin control UI
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/.openclaw/workspace/agent-workspace
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
